# backup/rsync/rbackup.service
[Unit]
Description=Local rsync backup to secondary drive
ConditionPathIsMountPoint=/mnt/local_backup

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle

# Pass display session for desktop notifications
PassEnvironment=DISPLAY WAYLAND_DISPLAY DBUS_SESSION_BUS_ADDRESS

ExecStart=%h/Scripts/backup/restic/rbackup.sh

# Prevent runaway resource usage
CPUQuota=50%
MemoryMax=512M

# Don't retry — the timer will catch the next cycle
Restart=no

# Timeout after 30 minutes
TimeoutStartSec=1800

[Install]
WantedBy=default.target
