# backup/restic/restic-backup.timer
[Unit]
Description=Run restic backup every 30 minutes

[Timer]
# Run 5 minutes after boot, then every 30 minutes
OnBootSec=5min
OnUnitActiveSec=30min

# Randomize within a 2-minute window to avoid thundering herd
RandomizedDelaySec=120

# Catch up on missed runs (e.g., laptop was asleep)
Persistent=true

[Install]
WantedBy=timers.target
